# etc/udev/rules.d/90-ux8406ma-keyboard.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="0b05", ATTR{idProduct}=="1b2c", MODE="0666"

ACTION=="add", SUBSYSTEM=="input", KERNEL=="event*", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="1b2c", ENV{ID_USB_INTERFACE_NUM}=="00", RUN+="/etc/ux8406ma/keyboard-connect-action.sh"

ACTION=="remove", SUBSYSTEM=="input", KERNEL=="event*", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="1b2c", ENV{ID_USB_INTERFACE_NUM}=="00", RUN+="/etc/ux8406ma/keyboard-disconnect-action.sh"